# restic/.config/systemd/user/restic-prune.timer
[Unit]
Description=Prune data from the restic repository weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
